# llama-factory/s6/services.d/llama-factory/run
#!/command/with-contenv bash

cd "${HOME}"
source .bashrc
echo "INFO: starting LLaMA Board..."
exec 2>&1
exec llamafactory-cli webui
